# init.d/pleroma
#!/sbin/openrc-run
supervisor=supervise-daemon
healthcheck_delay=60
healthcheck_timer=30
respawn_delay=5
respawn_max=0 # requires unreleased commit 54780a4 due to a bug
# change to 9001 or something if that's no good

name='Pleroma Server'
command="$pleroma_prefix/bin/pleroma"
command_args=start
command_user=pleroma:nobody
directory="$pleroma_prefix"

depend() {
    use postgresql
}

healthcheck() {
    [ -n "$pleroma_health" ] || return 0
    [ -x /usr/bin/curl ] || return 0
    /usr/bin/curl -q "localhost:${pleroma_port:-4000}/api/pleroma/healthcheck"
}
