# init.d/tcpcat4
#!/sbin/openrc-run
supervisor=supervise-daemon
healthcheck_timer=5
respawn_delay=30
respawn_period=300

: "${PORT:=${RC_SVCNAME#*.}}"
: "${DPORT:=$PORT}"

name=tcpcat4
command=socat
command_args="-4 TCP-LISTEN:$PORT,reuseaddr,fork,su=nobody TCP:$DADDR:$DPORT"
pidfile="/run/$RC_SVCNAME.pid"

extra_commands="checkconfig healthcheck"

depend() {
    need net
    after firewall
}

checkconfig() {
    # make sure all the variables are sane
    test "$PORT"  = "$((PORT))"  || return $?
    test "$DPORT" = "$((DPORT))" || return $?
    healthcheck
}

healthcheck() {
    nc -z "$DADDR" "$DPORT"
}

start_pre() {
    checkconfig
}
